# package/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop dynip-dns.service || true
systemctl disable dynip-dns.service || true
